# internal/cron.cron
*/2 * * * * /usr/bin/GET 'http://raks.com.ua/internal/advertise.php'  >> /var/www/raks/log/cron_mail.log
